# deploy/decard.service
[Unit]
Description=Decard API Server
After=network.target

[Service]
Type=simple
User=decard
Group=decard
WorkingDirectory=/opt/decard/back
Environment="PATH=/opt/decard/back/.venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/opt/decard/back/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
